# incident-notifier/incident-notifier-crash-alert.service
[Unit]
Description=Incident Notifier Crash Alert

[Service]
Type=oneshot
User=incident-notifier
Group=incident-notifier
EnvironmentFile=/etc/incident-notifier/secrets.env
ExecStart=/opt/incident-notifier/venv/bin/python /opt/incident-notifier/crash-alert.py %i
